# lab/orchestrator/demo-orchestrator.service
[Unit]
Description=Lab container orchestrator
After=docker.service
Requires=docker.service

[Service]
User=ubuntu
WorkingDirectory=/opt/demo-orchestrator
ExecStart=/usr/bin/node server.js
Restart=on-failure
RestartSec=3
Environment=PORT=8080
# Targets are defined in the challenge registry (challenges.js): each entry sets
# its own image, port, memory and success check. DEFAULT_CHALLENGE selects which
# one a session gets when it doesn't request a specific one (?challenge=<id>).
# Every challenge image must be pulled on the box at boot (user_data.sh.tftpl),
# since per-session containers run on a no-egress network.
Environment=DEFAULT_CHALLENGE=juice-admin
Environment=CLIENT_IMAGE=lab-client:latest
# One session = two containers. The Juice Shop targets need headroom, so this
# runs on a larger instance (see instance sizing in terraform locals); keep at 1.
Environment=MAX_SESSIONS=1
# Fallback target memory when a challenge omits memMb (challenges set their own).
Environment=TARGET_MEM_MB=1024
Environment=CLIENT_MEM_MB=128
# Guidance agent (Phase 4): non-secret config lives here; the secret BEDROCK_API_KEY
# is supplied by the EnvironmentFile below, written 0600 at boot (user_data.sh.tftpl).
# The leading "-" makes the file optional: with no key the orchestrator simply
# disables guidance (the hint endpoint 503s and the UI hides the control).
Environment=GUIDANCE_MODEL=google.gemma-4-31b
Environment=BEDROCK_REGION=us-west-2
Environment=MAX_CHAT_TURNS=20
EnvironmentFile=-/etc/demo-orchestrator.env

[Install]
WantedBy=multi-user.target
